# ArchersWebb/archerswebb.service
[Unit]
Description=ArchersWebb unified game deployment
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/archerswebb
Environment=NODE_ENV=production
Environment=WEB_HOST=0.0.0.0
Environment=WEB_PORT=8080
Environment=SHUTDOWN_TIMEOUT_MS=10000
ExecStart=/usr/bin/node /opt/archerswebb/start.js
Restart=always
RestartSec=3
User=www-data
Group=www-data
KillSignal=SIGTERM
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
